# weaviate_service/supervisor/services.conf
[supervisord]
nodaemon=true

[program:run_fastapi]
command=nohup /usr/local/bin/python3.9 /app/fastapi_runner.py
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true